# templates/teleclaude-docs-watch.path
[Unit]
Description=TeleClaude docs watch ({{PROJECT_ROOT}})

[Path]
PathModified={{PROJECT_ROOT}}/.agents
PathModified={{PROJECT_ROOT}}/docs/project
PathModified={{PROJECT_ROOT}}/docs/global
PathModified={{PROJECT_ROOT}}/teleclaude.yml
Unit={{UNIT_ID}}.service

[Install]
WantedBy=default.target
