# templates/teleclaude-docs-watch.service
[Unit]
Description=TeleClaude docs sync ({{PROJECT_ROOT}})

[Service]
Type=oneshot
WorkingDirectory={{PROJECT_ROOT}}
ExecStart=/bin/bash -lc '{{COMMAND}}'
